# scraper/crontab
# Run scraper daily at 2 AM
0 2 * * * root cd /app && python main.py >> logs/cron.log 2>&1

# Empty line required at end of cron file