# scraper/crontab
0 2 * * * cd /app && python main.py >> /app/logs/cron.log 2>&1
